# templates/thin.service.erb
[Unit]
Description=A fast and very simple Ruby web server
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/thin start --all /etc/thin/
ExecReload=/user/local/bin/thin restart --all /etc/thin/
ExecStop=/usr/local/bin/thin stop --all /etc/thin

[Install]
WantedBy=multi-user.target
